# keystat.service
[Unit]
Description=keystat capture service

[Service]
ExecStart=/usr/bin/keystat-capture /var/lib/keystat/map /dev/input/by-path/platform-i8042-serio-0-event-kbd

User=keystat
Group=keystat
SupplementaryGroups=input

IOSchedulingClass=realtime
IOSchedulingPriority=0
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99

CapabilityBoundingSet=
PrivateTmp=true
PrivateNetwork=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

StandardError=syslog

[Install]
WantedBy=multi-user.target
