# keystat.service
[Unit]
Description=keystat capture service

[Service]
ExecStart=/usr/bin/keystat-capture /var/lib/keystat/map /dev/input/by-path/platform-i8042-serio-0-event-kbd

CapabilityBoundingSet=dac_override

CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
IOSchedulingClass=realtime
IOSchedulingPriority=0

StandardError=syslog

[Install]
WantedBy=multi-user.target
